# speedtest.cron
*/10 * * * * zabbix /etc/zabbix/script/speedtest.sh -a
